# deploy/mywebsite-backend.service
[Unit]
Description=MyWebsite Flask backend (Gunicorn)
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/MyWebsite/backend
Environment=APP_DEBUG=0
Environment=PYTHONUNBUFFERED=1
ExecStart=/var/www/MyWebsite/backend/.venv/bin/gunicorn -c /var/www/MyWebsite/deploy/gunicorn.conf.py wsgi:app
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
